# examples/flask/supervisord.conf
[supervisord]
nodaemon=true

[program:uwsgi]
user=app
directory=/home/app
command=/usr/local/bin/uwsgi --http :8080 --pypy-wsgi-file app.py --processes 4 --threads 2
