# examples/flask/supervisord.conf
[supervisord]
nodaemon=true

[program:uwsgi]
user=root
directory=/home/app
command=/usr/local/bin/uwsgi --ini /home/app/uwsgi.conf
